# images/php-laravel/supervisord.conf
[supervisord]
nodaemon=true
logfile=/tmp/supervisord.log

[program:apache2]
command=/usr/sbin/apache2ctl -D FOREGROUND
autorestart=true
user=root

[program:php-fpm]
command=php-fpm
autorestart=true
