# images/php-laravel/supervisord.conf
[supervisord]
nodaemon=true

[program:apache2]
command=/usr/sbin/apache2ctl -D FOREGROUND
autorestart=true

[program:php-fpm]
command=php-fpm
autorestart=true
