# images/php-laravel/Alpine/supervisord.conf
[supervisord]
nodaemon=true
logfile=/tmp/supervisord.log

[program:httpd]
command=httpd -D FOREGROUND
autorestart=true

[program:php-fpm]
command=php-fpm
autorestart=true
